# udev_rules/90-netifchanges.rules
SUBSYSTEM=="net", RUN+="/opt/terminator/scripts/launcher.sh &"

